# services/sysvinit/inittab
zz:12345:respawn:/home/pi/projects/basketpinfo/sh/basket_start.sh
